# debian/web/postinst
#!/bin/sh
# postinst script for fossology-web
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
    # link the apache config for fossology
    if [ -f "/etc/apache2/sites-enabled/fossology.conf" ]; then
        rm -rf  /etc/apache2/sites-enabled/fossology.conf
    fi
    if [ -f "/etc/apache2/sites-available/fossology.conf" ]; then
        rm -rf  /etc/apache2/sites-available/fossology.conf
    fi
    PHP_PATH=$(php --ini | awk '/\/etc\/php.*\/cli$/{print $5}')
    phpIni="${PHP_PATH}/../apache2/php.ini"
    cp ${phpIni} ${phpIni}.bak
    sed -i 's/^\(max_execution_time\s*=\s*\).*$/\1300/' ${phpIni}
    sed -i 's/^\(memory_limit\s*=\s*\).*$/\1702M/' ${phpIni}
    sed -i 's/^\(post_max_size\s*=\s*\).*$/\1701M/' ${phpIni}
    sed -i 's/^\(upload_max_filesize\s*=\s*\).*$/\1700M/' ${phpIni}
    ln -s /etc/fossology/conf/fo-apache.conf /etc/apache2/sites-available/fossology.conf
    ln -s /etc/apache2/sites-available/fossology.conf /etc/apache2/sites-enabled/fossology.conf
    # need to run fo-postinstall web and agent stuff
    /usr/lib/fossology/fo-postinstall --web-only
    invoke-rc.d apache2 restart
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0


